# v18.04/overlay/etc/cron.d/qnap
*/1  *  *  *  * www-data /bin/bash -c "occ qnap:check-active-users"
